# utils/isard-scripts-update.service
[Unit]
After=network.target

[Service]
WorkingDirectory=<PATH>
ExecStart=sudo <PATH>/update.sh
User=<EXECUSER>

[Install]
WantedBy=default.target
