# utils/isard-scripts-update.service
[Unit]
After=network.target

[Service]
WorkingDirectory=<PATH>
ExecStart=<PATH>/update.sh

[Install]
WantedBy=default.target
