# utils/isard-scripts-first-run.service
[Unit]
After=network.target

[Service]
ExecStart=<PATH>/isard-scripts/run.sh

[Install]
WantedBy=default.target
