# utils/isard-scripts-update.service
[Unit]
After=network.target

[Service]
ExecStart=<PATH>/isard-scripts/update.sh

[Install]
WantedBy=default.target