# utils/isard-scripts-first-run.service
[Unit]
After=network.target

[Service]
WorkingDirectory=<PATH>
ExecStart=<PATH>/run.sh
User=<EXECUSER>

[Install]
WantedBy=default.target
