# utils/rc.local
#!/bin/bash

bash /etc/isard-scripts/run.sh
exit 0